# agent/docker/glastopf/dist/supervisor/glastopf.conf
[program:glastopf]
command= /opt/glastopf/glastopf-runner
autorestart=true